# deploy/xppm-dss.service
[Unit]
Description=XPPM Decision Support System
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/xppm-project
Environment="PATH=/home/ubuntu/xppm-project/.venv/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Comando de ejecución
ExecStart=/home/ubuntu/xppm-project/.venv/bin/python \
          /home/ubuntu/xppm-project/policy_server.py \
          --bundle /home/ubuntu/xppm-project/artifacts/deploy/v1 \
          --port 8000 \
          --host 0.0.0.0

# Restart policy
Restart=always
RestartSec=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xppm-dss

# Resource limits (opcional)
LimitNOFILE=65536
MemoryLimit=2G

[Install]
WantedBy=multi-user.target
